# docker-aaiteh.service.in
[Unit]
Description=Auto add ip address of new container in hosts file
After=docker.socket
BindsTo=docker.socket

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/docker-aaiteh.py

[Install]
WantedBy=multi-user.target
